# system-units-examples/gunicorn.socket
[Unit]
Description=gunicorn giftsBot socket

[Socket]
ListenStream=/run/gb_gunicorn.sock

[Install]
WantedBy=sockets.target
